# sv/openlist/run
#!/system/bin/sh
exec 2>&1

[ -r ./conf ] && . ./conf

exec chpst -u "${CHPST_USER:-shell:inet:sdcard_r:sdcard_rw}" \
    openlist server ${OPENLIST_ARGS:---data /storage/emulated/0/Android/openlist}
